# wbm-81c5/overlay/etc/crontabs/root
*/10 * * * *	lua /usr/bin/wbm-fwupdate
0 * * * * /usr/sbin/altermap-agent > /tmp/log/altermap-agent.log
