# wbm-6257/overlay/etc/crontabs/root
0 * * * * /usr/sbin/altermap-agent > /tmp/log/altermap-agent.log
