# wbm-5aa3/overlay/usr/lib/opkg/info/tinc.list
/etc/init.d/tinc
/usr/sbin/tincd
/lib/upgrade/keep.d/tinc
/etc/config/tinc
